# system-files/systemd/system/getty@tty1.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- shurtu-gal' --noclear --skip-login - $TERM
